# templates/systemd.service.j2
[Unit]
Description=Diaspora*

[Service]
ExecStart=/bin/bash -ec '{{ diaspora_home }}/start_server.sh'

[Install]
WantedBy=multi-user.target
