# templates/systemd.service.j2
[Unit]
Description=Diaspora*

[Service]
ExecStart=/bin/bash -ec '/home/{{ os_user }}/start_server.sh'

[Install]
WantedBy=multi-user.target
